# deploy/property-tracker-watchdog.timer
[Unit]
Description=Run property tracker watchdog every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
